# deploy/trading.service
[Unit]
Description=srb-algo trading engine
After=network.target

[Service]
Type=simple
User=%i
# Systemd unit prefilled for Hostinger deploy user `u570718221`.
# Update if deploying under a different user.
WorkingDirectory=/home/u570718221/domains/srbalgo.bytevantage.in/public_html/srb-algo
EnvironmentFile=/home/u570718221/domains/srbalgo.bytevantage.in/public_html/srb-algo/.env
ExecStart=/home/u570718221/domains/srbalgo.bytevantage.in/public_html/srb-algo/venv/bin/python3 backend/main.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=srb-algo

[Install]
WantedBy=multi-user.target
